# deploy/exo.service
[Unit]
Description=Exokephalos
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/erwin/.local/bin/xo serve
Environment=EXO_DIR=/home/erwin/exo
Restart=on-failure
RestartSec=5

# Filesystem hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/erwin/exo

[Install]
WantedBy=default.target
